# .docker/crontab/cronjob_24hour.sh
#!/bin/sh
echo "$(date): executed script - 24 hour"
