# .docker/crontab/cronjob_weekly.sh
#!/bin/sh
echo "$(date): executed script - weekly"
