# .docker/crontab/cronjob_1hour.sh
#!/bin/sh
echo "$(date): executed script - 1 hour"
